# aliada/aliada-links-discovery/src/test/resources/aliada_links_discovery.cron
29 8 8 12 * 6104978122374998935/links-discovery-task-runner.sh -923546869 164360400 -3559353425198174491
